# support/packaging/common/mesos-agent.service.in
[Unit]
Description=ClusterD Agent
After=network.target
Wants=network.target

[Service]
ExecStart=@RUN_PATH@/bin/mesos-init-wrapper agent
KillMode=process
Restart=always
RestartSec=20
LimitNOFILE=16384
CPUAccounting=true
MemoryAccounting=true
Delegate=true

[Install]
WantedBy=multi-user.target
