# support/packaging/common/mesos-master.service.in
[Unit]
Description=ClusterD Master
After=network.target
Wants=network.target

[Service]
ExecStart=@RUN_PATH@/bin/mesos-init-wrapper master
Restart=always
RestartSec=20
LimitNOFILE=16384

[Install]
WantedBy=multi-user.target
